# demo/odoo_demo.conf
[program:odoo-demo]
<<<<<<< HEAD
command=pew in odoo-7.0 /usr/local/lib/odoo/openerp-server -d demo_odoo70 --addons-path=addons
=======
command=/usr/local/lib/odoo/openerp-server -d demo_odoo70 --addons-path=/usr/local/lib/odoo/addons
>>>>>>> 56bbbcdfee685ed28964c6eca12f2ba5d319ab63
user=openerp
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
environment=HOME="/opt/openerp",PATH="{VENV_PATH}/bin:{PATH}",VIRTUAL_ENV="{VENV_PATH}"
priority=50
