# demo/odoo_demo.conf
[program:odoo-demo]
command="/usr/local/lib/openerp/openerp-server --addons-path=web,addons -d odoo_demo70"
user=openerp
autorestart=true